# ops/spectrum-classifier-health/spectrum-classifier-health.service
[Unit]
Description=rf_luv classifier health monitor (one-shot)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3.11 %h/dev/rf_luv/spectrum/classifier_health.py
Environment=CLICKHOUSE_HOST=localhost
Environment=CLICKHOUSE_PORT=8126
Environment=CLICKHOUSE_DB=spectrum
Environment=CLICKHOUSE_USER=spectrum
Environment=CLICKHOUSE_PASSWORD=spectrum_local
TimeoutStartSec=60
Restart=no

[Install]
WantedBy=default.target
