# ops/rtl-tcp/rtl-tcp-watchdog.service
[Unit]
Description=Probe rtl_tcp; recover on sample starvation
After=rtl_tcp.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rtl-tcp-watchdog.py
