# ops/rtl-tcp/rtl-reset-failed.service
[Unit]
Description=Clear systemd failure state on rtl_tcp/rtl-scanner user units

[Service]
Type=oneshot
# Clear "failed" state on any rtl-* template or singleton unit so that
# Restart=always can re-engage. Prevents the 2026-04-23 pathology where
# rtl-tcp@v3-01 hit StartLimitBurst=20 in a transient USB flap and stayed
# permanent-failed until manual intervention. Zero-cost when everything is
# healthy — reset-failed is a no-op on non-failed units.
#
# Arguments are globs only — literal singleton names (rtl_tcp.service,
# rtl-tcp-watchdog.service) would cause systemctl to exit 1 with "Unit not
# loaded" after those singletons were masked to /dev/null on 2026-04-23.
# Globs against never-loaded templates are no-ops (exit 0), so the set
# below covers live templates plus self without ever failing.
ExecStart=/bin/systemctl --user reset-failed rtl-tcp@*.service rtl-scanner@*.service rtl-tcp-watchdog@*.service rtl-reset-failed.service
