# ops/rtl-tcp/rtl-tcp-watchdog.timer
[Unit]
Description=rtl_tcp health probe (every 30s)

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
