# ops/spectrum-monitor/freshness-probe.service
[Unit]
Description=Spectrum pipeline freshness probe (ClickHouse-level)
Wants=network.target
After=network.target

[Service]
Type=oneshot
# Optional env file for runtime overrides (EXPECTED_DONGLES, FRESHNESS_WARN_S,
# FRESHNESS_CRITICAL_S, etc.). The leading `-` makes the file optional —
# defaults baked into the script are used if the file is missing.
EnvironmentFile=-/etc/rtl-scanner/freshness-probe.env
ExecStart=/usr/local/bin/rf-freshness-probe
TimeoutStartSec=60
User=root
