# ops/notify/heartbeat.service
[Unit]
Description=Daily ntfy heartbeat — confirms alert pipe is alive
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rf-heartbeat
TimeoutStartSec=30
