# ops/spectrum-acars-feedback/spectrum-acars-feedback.service
[Unit]
Description=Promote ACARS decode confirmations into spectrum.listening_log
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3.11 %h/dev/rf_luv/spectrum/acars_feedback.py
# Spectrum side
Environment=CLICKHOUSE_HOST=localhost
Environment=CLICKHOUSE_PORT=8126
Environment=CLICKHOUSE_DB=spectrum
Environment=CLICKHOUSE_USER=spectrum
Environment=CLICKHOUSE_PASSWORD=spectrum_local
# ACARS side (separate ClickHouse instance — see acars/docker-compose.yml)
Environment=ACARS_CLICKHOUSE_HOST=localhost
Environment=ACARS_CLICKHOUSE_PORT=8127
Environment=ACARS_CLICKHOUSE_DB=acars
Environment=ACARS_CLICKHOUSE_USER=acars
Environment=ACARS_CLICKHOUSE_PASSWORD=acars_local
# Threshold tuning — adjust if LGAV ACARS volume shifts.
Environment=ACARS_FEEDBACK_LOOKBACK_HOURS=24
Environment=ACARS_FEEDBACK_MIN_MESSAGES=10
TimeoutStartSec=60
Restart=no

[Install]
WantedBy=default.target
